# supervisor-happywed-worker.conf
[program:happywed-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/your/happywed/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
directory=/path/to/your/happywed
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/path/to/your/happywed/storage/logs/worker.log
stopwaitsecs=3600
